# roles/user/files/.config/autostart/clipit-startup.desktop
[Desktop Entry]
Name=ClipIt
Icon=gtk-paste
Exec=clipit
Terminal=false
Type=Application
